# supervisord-catalina.conf
[program:catalina]
command=/opt/bonita/BonitaBPMCommunity-7.4.1-Tomcat-7.0.67/server/bin/catalina.sh run
numprocs=1
autostart=true
autorestart=true